# inittab
::respawn:/usr/local/bin/start_syslogd
::respawn:crond -f
::respawn:/usr/local/bin/start_named
::respawn:/usr/local/bin/start_root_proxy 1
::respawn:/usr/local/bin/start_root_proxy 2
::respawn:/usr/local/bin/rest.mjs
::respawn:/usr/local/bin/start_eth_redirect
::respawn:/usr/local/bin/start_exporter
::respawn:/usr/local/bin/start_dnsdist
::respawn:/usr/sbin/haproxy -f /usr/local/etc/haproxy.cfg
